# ubiquiti/UAP-AC-M/usr/lib/opkg/info/kmod-qca-wifi-unified-perf.list
/lib/modules/4.4.153/asf.ko
/lib/modules/4.4.153/qca_da.ko
/usr/sbin/thermaltool
/lib/modules/4.4.153/ath_dfs.ko
/lib/modules/4.4.153/mem_manager.ko
/lib/modules/4.4.153/ath_spectral.ko
/etc/init.d/ssid_steering
/usr/sbin/apstats
/usr/sbin/regdomain
/lib/modules/4.4.153/ath_dev.ko
/usr/sbin/spectraltool
/usr/sbin/athchans
/lib/modules/4.4.153/ath_rate_atheros.ko
/etc/config/ssid-steering
/usr/sbin/80211stats
/lib/modules/4.4.153/ath_hal.ko
/usr/sbin/ol_regdomain
/usr/sbin/radartool
/lib/wifi/qcawifi_countrycode.txt
/lib/modules/4.4.153/qca_ol.ko
/usr/sbin/wifitool
/usr/sbin/athstats
/lib/modules/4.4.153/qdf.ko
/lib/wifi/qca-wifi-modules
/lib/sbin/wifi_try
/usr/sbin/wlanconfig
/lib/modules/4.4.153/ath_pktlog.ko
/lib/modules/4.4.153/umac.ko
/usr/sbin/athstatsclr
/lib/wifi/qcawifi.sh
/usr/sbin/wifi_try
