# ubiquiti/UAP-AC-M/usr/lib/opkg/info/ubnt-snmpd.list
/etc/init.d/snmpd
/etc/config/snmpd
/etc/snmp/snmp.conf
/usr/sbin/snmpd
