# ubiquiti/U6-IW/etc/init.d/ubnt_plugins
#!/bin/sh /etc/rc.common
# vim: ft=sh

START=12
STOP=10

exec 2> /tmp/rc.txt
exec 1>&2
set -x

start() {
    . /usr/etc/rc.d/rc start
}

stop() {
    . /usr/etc/rc.d/rc stop
}

