# tplink_en/Archer AX23/usr/lib/opkg/info/agile_config.list
/sbin/agile_config.sh
/sbin/agileconfig-countdown-1h
/etc/init.d/agile_config
/sbin/reset
/sbin/agile_config_calc_1h.sh
