# tplink_en/TL-WR941HP/usr/lib/opkg/info/kmod-qca-wifi-perf.list
/lib/modules/3.3.8/asf.ko
/lib/wifi/tplink.sh.orig
/lib/wifi/tplink.sh
/lib/modules/3.3.8/ath_hal.ko
/lib/modules/3.3.8/ath_rate_atheros.ko
/lib/modules/3.3.8/ath_pktlog.ko
/lib/modules/3.3.8/adf.ko
/lib/debug/qca-wifi
/etc/init.d/art
/lib/wifi/tplink-update-uci
/etc/modules.d/33-qca-wifi-perf
/etc/init.d/wireless
/lib/wifi/qcawifi_countrycode.txt
/usr/sbin/tx99tool
/lib/modules/3.3.8/hst_tx99.ko
/sbin/wifi_try
/lib/modules/3.3.8/umac.ko
/usr/sbin/wlanconfig
/lib/modules/3.3.8/ath_dev.ko
/sbin/is_cal
/etc/hotplug.d/iface/02-wireless-vlan
/lib/wifi/qcawifi.sh
