# tplink_en/TL-WR941HP/usr/lib/opkg/info/operation-mode.list
/etc/init.d/operation_mode
/etc/disable_module_re
/etc/disable_list_re
/etc/init.d/operation_module
/etc/disable_module_ap
/etc/disable_list_ap
