# tplink_en/TL-WR941HP/usr/lib/opkg/info/switch-setup.list
/lib/switch/config.sh
/etc/hotplug.d/iface/90-switch_setup
/lib/switch/core.sh
/etc/init.d/switch
/etc/init.d/phyport
/etc/config/switch
