# tplink_en/Archer C2600/etc/hotplug.d/firewall/02-wireless-vlan
[ "$ACTION" == "start" ] && {
    /sbin/wifi vlan
}

[ "$ACTION" == "stop" ] && {
    : # Null command
}
