# tplink_en/Archer C59/etc/init.d/shift-check-start
#!/bin/sh /etc/rc.common

START=13
SERVICE_DAEMONIZE=1
cycle()
{
	sleep 2 
	ledcli SHIFT_CHECK
}
start()
{
	cycle &
      
}
