# netgear/RS700S/usr/bin/circlev2/shares/usr/bin/service/apid/run
#!/bin/sh
exec $CIRCLE_ROOT/apid -u root -E $CIRCLE_DATA/apid.pem  -D
