# netgear/RAX29/etc/init.d/create_cronjob.init
#!/bin/sh /etc/rc.common
START=99

start() {
    /usr/bin/seal_randcron daily_fing_query
    /usr/bin/seal_randcron daily_spc_update
    /etc/init.d/cron reload
}

boot() {
    start
}
