# tplink_en/Archer C59/usr/lib/opkg/info/kmod-qca-wifi-10.4-perf.list
/lib/modules/3.3.8/asf.ko
/lib/wifi/tplink.sh
/lib/modules/3.3.8/ath_hal.ko
/etc/hotplug.d/button/40-wireless
/lib/modules/3.3.8/ath_rate_atheros.ko
/lib/modules/3.3.8/ath_pktlog.ko
/lib/modules/3.3.8/adf.ko
/etc/modules.d/33-qca-wifi-10.4-perf
/lib/debug/qca-wifi
/lib/wifi/tplink-update-uci
/etc/init.d/wireless
/lib/wifi/qcawifi_countrycode.txt
/usr/lib/libtlvparser.so
/usr/lib/libtlvencoder.so
/usr/sbin/athtestcmd
/usr/sbin/tx99tool
/usr/lib/libtlvtemplate.so
/lib/modules/3.3.8/hst_tx99.ko
/usr/lib/libtlvcmdrsp.so
/sbin/wifi_try
/lib/modules/3.3.8/umac.ko
/lib/sbin/wifi_try
/usr/lib/libtlvutil.so
/usr/sbin/wlanconfig
/etc/init.d/wireless_sta
/usr/sbin/athdiag
/lib/modules/3.3.8/ath_dev.ko
/etc/init.d/zwireless_check
/sbin/is_cal
/etc/hotplug.d/iface/02-wireless-vlan
/lib/wifi/qcawifi.sh
