# tplink_en/Archer C59/usr/lib/opkg/info/ledctrl-apps.list
/usr/bin/wps_led
/usr/sbin/multi_lan_connected
/usr/sbin/online-test
/etc/init.d/switch_led
/usr/sbin/lan_connected
/usr/bin/switch_led
/usr/sbin/wan_connected
