# tplink_en/Archer C59/usr/lib/opkg/info/tddp.list
/usr/bin/tddp
/etc/init.d/tddp
/sbin/tddp_hook_erase_radio
