# tplink_en/Archer C59/usr/lib/opkg/info/usbshare.list
/etc/hotplug.d/iface/70-ftpex
/sbin/usbuser
/sbin/ftpex
/etc/config/sharecfg
/sbin/ftp_alias
/etc/init.d/usbshare
/usr/sbin/usbpoweron
/etc/config/usbshare
/usr/sbin/usbshare
