# tplink_en/Archer AX95/usr/lib/opkg/info/access-control.list
/lib/access_control/core.sh
/lib/access_control/core_global.sh
/lib/access_control/core_log.sh
/sbin/ac_macfilter_wrap_black
/lib/upgrade/keep.d/access-control
/etc/init.d/access_control
/lib/access_control/config.sh
/sbin/ac
/sbin/ac_macfilter_wrap_white
/etc/hotplug.d/firewall/20-access_control
/etc/config/access_control
