# tplink_en/Archer AX95/usr/lib/opkg/info/client-speed-limit.list
/etc/init.d/client_speed_limit
/etc/hotplug.d/firewall/97-client_speed_limit
/lib/client_speed_limit/csl_core.sh
/etc/config/client_speed_limit
/lib/upgrade/keep.d/client-speed-limit
/lib/client_speed_limit/csl_ratelimit.sh
/lib/client_speed_limit/csl_config.sh
/etc/hotplug.d/iface/97-client_speed_limit
