# tplink_en/Archer AX95/usr/lib/opkg/info/client_mgmt.list
/usr/bin/client_mgmt
/usr/sbin/new_client_access_report
/etc/init.d/client_mgmt
