# tplink_en/Archer AX95/usr/lib/opkg/info/cloud-report.list
/usr/sbin/report_upload_url_apply
/usr/sbin/report_upload_time_apply
/usr/sbin/report_upload_components
/usr/sbin/report_upload_period
/usr/lib/lua/cloud_report/report_error.lua
/usr/sbin/report_upload_unbind
/etc/init.d/cloud_report
/usr/sbin/report_upload_alert
/etc/init.d/avira_get_service_status_from_cloud
/usr/sbin/report_get_device_status
/usr/sbin/report_upload_bind
/usr/sbin/cloud-report
/usr/sbin/report_upload_daily
/usr/sbin/report_download_elapsedTime
