# tplink_en/Archer AX95/usr/lib/opkg/info/guestnetwork-effectivetime-ctrl.list
/etc/init.d/guestnetwork_effectivetime_ctrl
/lib/upgrade/keep.d/guestnetwork-effectivetime-ctrl
/lib/guestnetwork_effectivetime_ctrl/core.sh
/sbin/gec_check
/lib/guestnetwork_effectivetime_ctrl/core_gec.sh
/etc/config/guestnetwork_effectivetime_ctrl
