# tplink_en/Archer AX95/usr/lib/opkg/info/guestnetwork-bandwidth-ctrl.list
/etc/config/guestnetwork_bandwidth_ctrl
/lib/guestnetwork_bandwidth_ctrl/gbc_ratelimit.sh
/lib/guestnetwork_bandwidth_ctrl/gbc_config.sh
/etc/hotplug.d/firewall/98-guestnetwork_bandwidth_ctrl
/lib/guestnetwork_bandwidth_ctrl/gbc_core.sh
/etc/hotplug.d/iface/98-guestnetwork_bandwidth_ctrl
/lib/upgrade/keep.d/guestnetwork-bandwidth-ctrl
/etc/init.d/guestnetwork_bandwidth_ctrl
