# tplink_en/Archer AX95/usr/lib/opkg/info/safethings.list
/usr/share/.avira/avirasentinelfull/avirasentinelfull_220110_FB-1_prod/passwords.txt
/usr/lib/avira/libstdc++.so.6
/usr/bin/avira/aviraserviceselector/aviraserviceselector_220110_FB-1_prod/uninstallComplete.sh
/usr/lib/avira/libz.txt
/usr/lib/avira/libstdc++.so.6.0.27
/usr/lib/avira/libm-2.30.so
/usr/share/.avira/avirasentinelfull/avirasentinelfull_220110_FB-1_prod/ports.csv
/usr/lib/avira/librt.txt
/usr/lib/avira/libcap.so.2.27
/usr/lib/avira/libnetfilter_conntrack.so.3
/usr/bin/avira/avirasentinelfull/avirasentinelfull_220110_FB-1_prod/avirasentinelfull.sh
/usr/lib/avira/libmbedx509.so.2.16.7
/usr/bin/avira/avirasentinellite/avirasentinellite_220110_FB-1_prod/avirasentinellite.sh
/usr/lib/avira/libc.txt
/usr/lib/avira/libz.so.1.2.11
/usr/lib/avira/libsodium.so.23
/usr/lib/avira/libmbedtls.so.2.16.7
/usr/lib/avira/libnetfilter_conntrack.txt
/usr/bin/avira/aviraserviceselector/aviraserviceselector
/usr/share/.avira/aviraserviceselector/aviraserviceselector_220110_FB-1_prod/mainConfig.config.sig
/usr/share/.avira/ca-certificates.config.sig
/usr/lib/avira/libsodium.txt
/usr/lib/avira/libm.so.6
/usr/lib/avira/libwebsockets.txt
/usr/lib/avira/libgcc_s.so.1
/usr/bin/avira/avirasentinellite/avirasentinellite_220110_FB-1_prod/avirasentinellite
/usr/bin/avira/aviraserviceselector/aviraserviceselector_220110_FB-1_prod/aviraserviceselector
/usr/lib/avira/libz.so.1
/usr/lib/avira/libsodium.so.23.3.0
/usr/lib/avira/libnfnetlink.so.0
/usr/share/.avira/avirasentinellite/avirasentinellite_220110_FB-1_prod/mainConfig.config
/usr/bin/avira/avirasentinellite/avirasentinellite_220110_FB-1_prod/uninstallSpecific.sh
/etc/init.d/aviraserviceselector
/usr/share/.avira/avirasentinelfull/avirasentinelfull_220110_FB-1_prod/mainConfig.config.sig
/etc/init.d/avirasentinelfull
/etc/init.d/aviraservicemaster
/usr/lib/avira/libmbedtls.txt
/usr/lib/avira/libm.txt
/usr/lib/avira/libmbedcrypto.txt
/usr/lib/avira/libstdc++.txt
/etc/init.d/avirasentinellite
/usr/lib/avira/libcap.txt
/usr/bin/avira/avirasentinelfull/avirasentinelfull_220110_FB-1_prod/uninstallComplete.sh
/usr/bin/avira/aviraserviceselector/aviraserviceselector_220110_FB-1_prod/uninstallSpecific.sh
/usr/share/.avira/avirasentinellite/avirasentinellite_220110_FB-1_prod/mainConfig.config.sig
/usr/lib/avira/libpcap.so.1
/usr/lib/avira/libpthread.so.0
/usr/lib/avira/libcap.so.2
/usr/lib/avira/libmnl.so.0
/usr/bin/avira/avirasentinelfull/avirasentinelfull
/usr/lib/avira/libc-2.30.so
/usr/lib/avira/libcurl.so.4.6.0
/usr/bin/avira/sshpass
/usr/lib/avira/libnfnetlink.txt
/usr/share/.avira/aviraserviceselector/aviraserviceselector_220110_FB-1_prod/mainConfig.config
/usr/share/.avira/avirasentinellite/avirasentinellite_220110_FB-1_prod/passwords.txt
/usr/sbin/aviraservice
/usr/lib/avira/libcurl.txt
/usr/lib/avira/libmbedcrypto.so.3
/usr/lib/avira/libmbedx509.so.0
/usr/lib/avira/libc.so.6
/usr/bin/avira/aviraserviceselector/aviraserviceselector_220110_FB-1_prod/aviraserviceselector.sh
/usr/lib/avira/libmnl.so.0.2.0
/usr/lib/avira/librt.so.1
/usr/lib/avira/libwebsockets.so.13
/usr/lib/avira/libnfnetlink.so.0.2.0
/usr/lib/avira/libcurl.so.4
/usr/lib/avira/libmbedx509.txt
/usr/lib/avira/libpthread.txt
/usr/lib/avira/libmbedcrypto.so.2.16.7
/usr/share/.avira/ca-certificates.config
/usr/bin/avira/avirasentinelfull/avirasentinelfull_220110_FB-1_prod/avirasentinelfull
/usr/bin/avira/avirasentinellite/avirasentinellite
/usr/lib/avira/libpthread-2.30.so
/usr/lib/avira/libgcc_s.txt
/usr/lib/avira/libmnl.txt
/usr/lib/avira/libpcap.so.1.9.1
/usr/lib/avira/librt-2.30.so
/usr/bin/avira/avirasentinellite/avirasentinellite_220110_FB-1_prod/uninstallComplete.sh
/usr/bin/avira/avirasentinelfull/avirasentinelfull_220110_FB-1_prod/uninstallSpecific.sh
/usr/lib/avira/libmbedtls.so.12
/usr/lib/avira/libpcap.txt
/usr/share/.avira/avirasentinellite/avirasentinellite_220110_FB-1_prod/ports.csv
/usr/share/.avira/avirasentinelfull/avirasentinelfull_220110_FB-1_prod/mainConfig.config
/usr/lib/avira/libnetfilter_conntrack.so.3.7.0
