# tplink_en/Archer AX95/usr/lib/opkg/info/snort.list
/etc/init.d/tp_security
/etc/snort/preproc_rules/preprocessor.rules
/etc/snort/gen-msg.map
/etc/snort/preproc_rules/sensitive-data.rules
/etc/snort/attribute_table.dtd
/etc/snort/reference.config
/etc/snort/unicode.map
/etc/snort/rules/dns.rules
/etc/snort/rules/penetration.rules
/etc/snort/preproc_rules/decoder.rules
/etc/snort/snort.conf
/etc/snort/rules/dos.rules
/usr/bin/snort
/etc/snort/classification.config
/etc/snort/rules/scan.rules
/etc/snort/rules/bpf_filter.rules
/etc/snort/rules/ddos.rules
/etc/config/snort
/etc/snort/rules/community-dos.rules
