# tplink_en/Archer AX95/usr/lib/opkg/info/watchdog.list
/usr/bin/avira/avirawatchdog/avirawatchdog_220110_FB-1_prod/avirawatchdog
/usr/bin/avira/avirawatchdog/avirawatchdog_220110_FB-1_prod/uninstallComplete.sh
/etc/init.d/avirawatchdog
/usr/bin/avira/avirawatchdog/avirawatchdog_220110_FB-1_prod/avirawatchdog.sh
/usr/share/.avira/avirawatchdog/avirawatchdog_220110_FB-1_prod/mainConfig.config
/usr/bin/avira/avirawatchdog/avirawatchdog_220110_FB-1_prod/uninstallSpecific.sh
/usr/bin/avira/avirawatchdog/avirawatchdog
/usr/scripts/avira/sendWatchdogHeartbeat.sh
/usr/share/.avira/avirawatchdog/avirawatchdog_220110_FB-1_prod/mainConfig.config.sig
/usr/scripts/avira/disableWatchdogHeartbeat.sh
/usr/scripts/avira/watchdogCronTask.sh
/usr/scripts/avira/enableWatchdogHeartbeat.sh
