# tplink_en/Archer AX95/usr/lib/opkg/info/wifidog.list
/usr/lib/libhttpd.so.0.0.0
/etc/wifidog/locale/es_ES/lan.js
/etc/wifidog/locale/da_DK/lan.js
/etc/wifidog/Enable2.png
/usr/lib/libhttpd.a
/usr/bin/wifidog
/etc/wifidog/locale/fr_FR/lan.js
/etc/wifidog/locale/uk_UA/lan.js
/etc/wifidog/locale/th_TH/lan.js
/etc/wifidog/locale/sk_SK/lan.js
/etc/wifidog/locale/pt_PT/lan.js
/etc/wifidog/locale/ko_KR/lan.js
/etc/wifidog/locale/es_MX/lan.js
/etc/hotplug.d/iface/99-wifidog
/etc/wifidog/locale/pt_BR/lan.js
/usr/lib/libhttpd.so
/etc/wifidog/locale/sv_SE/lan.js
/etc/wifidog/locale/no_NO/lan.js
/etc/hotplug.d/firewall/99-wifidog
/etc/wifidog/locale/fi_FI/lan.js
/etc/wifidog/portal_back.jpg
/etc/wifidog/portal.html
/etc/wifidog/lan.js
/etc/wifidog/locale/cs_CZ/lan.js
/etc/wifidog/jquery.min.js
/etc/wifidog/locale/jp_JP/lan.js
/etc/wifidog/locale/tr_TR/lan.js
/etc/wifidog/locale/zh_TW/lan.js
/etc/wifidog/locale/en_US/lan.js
/etc/wifidog/wifidog-msg.html
/usr/include/httpd.h
/etc/wifidog/Invisible.png
/etc/init.d/wifidog
/usr/lib/libhttpd.so.0
/etc/wifidog/portal_logo.png
/usr/bin/wdctl
/etc/wifidog/wifidog.conf
/etc/wifidog/locale/ro_RO/lan.js
/lib/preinit/73_init_www_ramfs
/etc/wifidog/locale/nl_NL/lan.js
/etc/wifidog/red-eclamation.png
/etc/wifidog/Visible.png
/etc/wifidog/Close.png
/etc/wifidog/locale/ru_RU/lan.js
/etc/wifidog/locale/bg_BG/lan.js
/etc/wifidog/locale/pl_PL/lan.js
/etc/wifidog/locale/hu_HU/lan.js
/etc/wifidog/locale/vi_VN/lan.js
/usr/lib/libhttpd.la
/etc/wifidog/Enable1.png
/etc/wifidog/locale/it_IT/lan.js
/etc/wifidog/locale/de_DE/lan.js
