# tplink_en/Archer AX95/usr/lib/opkg/info/wireless-script.list
/lib/wifi/triband/wlan.ini
/lib/wifi/dualband/wlan.ini
/sbin/tddp_hook_enable_testmode
/lib/wifi/wps_button
/lib/wifi/brcm_wlan_var.sh
/sbin/wifi_firm
/sbin/tddp_hook_erase_radio
/lib/wifi/tplink_brcm.sh
/etc/rc.d/S99zzzzzwifi_check
/lib/wifi/wl_iface
/lib/wifi/config.ini
/lib/wifi/brcm_wlan_extend.sh
