# tplink_en/AD7200/etc/hotplug.d/iface/01-wportal
#when lan ip changed, need to restart wportal, rmmod here .

exit 0

[ "$INTERFACE" == "lan" ] && wportalctrl -c && echo "stop" > /tmp/wportal/status