# tplink_en/Archer C58HP/usr/lib/opkg/info/autoreboot.list
/etc/init.d/auto_reboot
/lib/auto_reboot.sh
