# tplink_en/Archer C58HP/usr/lib/opkg/info/client_mgmt.list
/usr/bin/client_mgmt
/etc/init.d/client_mgmt
/etc/hotplug.d/iface/70-client_mgmt
