# tplink_en/Archer C58HP/usr/lib/opkg/info/diehard.list
/sbin/diehard_monitor
/etc/init.d/diehard
/lib/diehard/diehard_core.sh
