# tplink_en/Archer C58HP/usr/lib/opkg/info/kmod-qca-wifi-10.4-unified-perf.list
/lib/modules/3.3.8/asf.ko
/lib/wifi/tplink.sh
/lib/modules/3.3.8/ath_hal.ko
/lib/modules/3.3.8/ath_rate_atheros.ko
/lib/modules/3.3.8/ath_pktlog.ko
/lib/modules/3.3.8/adf.ko
/etc/modules.d/33-qca-wifi-10.4-unified-perf
/etc/init.d/art
/etc/init.d/ssid_steering
/lib/wifi/tplink-update-uci
/etc/config/ssid-steering
/etc/init.d/wireless
/lib/wifi/qcawifi_countrycode.txt
/usr/sbin/tx99tool
/lib/modules/3.3.8/hst_tx99.ko
/lib/modules/3.3.8/umac.ko
/lib/sbin/wifi_try
/usr/sbin/wlanconfig
/lib/modules/3.3.8/smart_antenna.ko
/lib/modules/3.3.8/ath_dev.ko
/sbin/is_cal
/etc/hotplug.d/iface/02-wireless-vlan
/usr/sbin/wifi_try
/lib/wifi/qcawifi.sh
