# tplink_en/Archer C58HP/usr/lib/opkg/info/miniupnpd.list
/usr/sbin/miniupnpd
/usr/sbin/upnp_conflict_check
/etc/config/upnpd
/etc/hotplug.d/iface/50-miniupnpd
/usr/share/miniupnpd/firewall.include
/etc/uci-defaults/99-miniupnpd
/etc/init.d/miniupnpd
