# tplink_en/Archer C58HP/usr/lib/opkg/info/repeater-monitor.list
/sbin/repeater_monitor
/etc/init.d/repeater
/lib/wifi/repeater_core.sh
