# tplink_en/Archer C58HP/usr/lib/opkg/info/rp-pppoe-relay.list
/usr/sbin/pppoe-relay
/etc/init.d/pppoe-relay
/etc/config/pppoe
/etc/default/pppoe-relay
/etc/hotplug.d/iface/65-pppoe
