# tplink_en/Archer C58HP/usr/lib/opkg/info/smartdhcp.list
/etc/init.d/smartdhcp
/etc/hotplug.d/smartdhcp/01-smartdhcp
/lib/smartdhcp/smartdhcp_core.sh
/lib/smartdhcp/smartdhcp.script
/sbin/smartdhcp_monitor
/lib/smartdhcp/smartdhcp_functions.sh
