# ubiquiti/U7-Pro/usr/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyMSM0::askfirst:/sbin/getty -L ttyMSM0 115200 vt100
null::respawn:/usr/bin/reset-handler
null::respawn:/bin/thermald_wrapper.sh
null::respawn:/bin/fanctrl_wrapper.sh
