# netgear/RAX38/etc/hotplug.d/iface/70-update-ftp
#!/bin/sh

if [ "$ACTION" = "ifup" ]; then
        /usr/bin/objcli modify updateFtp $DEVICE
fi
