# netgear/RS700S/usr/bin/circlev2/shares/usr/bin/service/iplistd/run
#!/bin/sh
exec $CIRCLE_ROOT/iplistd -i $CIRCLE_BRIDGE -d
