# ubiquiti/U7-Pro/etc/init.d/ubnt_bootconf
#!/bin/sh /etc/rc.common

START=11

boot() {
	# process actions that must be done after "boot" (module loading)
	. /usr/etc/rc.bootconf
}
