# tplink_en/Archer AX50/etc/hotplug.d/iface/00-zboard_intel
[ "internet" = "$INTERFACE" -o "wan" = "$INTERFACE" -o "lan" = "$INTERFACE" -o "lanv6" = "$INTERFACE"  ] && {
	case "$ACTION" in
	ifup)
		config_ppa
		;;
	ifdown)

		;;
	esac;
}