# tplink_en/Archer AX12/etc_ro/hotplug.d/net/00-guest-access-local
#!/bin/sh

[ "change" = "$ACTION" -a -n "$SUBACTION" -a -n "$BRIDGENAME" -a -n "$INTERFACE" ] && {
    case $SUBACTION in
        braddif)
            wifi vlan with_br "$BRIDGENAME" "$INTERFACE"
            ;;
        brdelif)
            ;;
    esac
}
