# tplink_en/Archer C58HP/etc/hotplug.d/iface/70-client_mgmt
#!/bin/sh

if [ "$ACTION" == "ifup" ] && [ "$DEVICE" == "br-lan" ] && [ "$INTERFACE" == "lan" ]; then
    ./etc/init.d/client_mgmt restart
fi

