# tplink_en/Archer C4000/usr/lib/opkg/info/powermanager.list
/etc/init.d/powermanager
/usr/sbin/powermanager
/etc/powermanager.conf
