# tplink_en/Archer C4000/usr/lib/opkg/info/speed-test.list
/etc/speed-test/settings.txt
/usr/bin/speed-test
/etc/init.d/speedtest
/usr/sbin/speedtest
/etc/speed-test/speed_test.sh
