# netgear/RAX5/etc/uci-defaults/32_gen_NTGR_WiFi
#!/bin/sh

/usr/bin/lua -lcommonFunc/wifiUtils_convertWifiSettings -e "local M = require \"commonFunc.wifiUtils_convertWifiSettings\"; M.update_NTGR_WiFi_uci()"
