# netgear/RAX5/etc/uci-defaults/netgear.default
#!/bin/sh
. /usr/share/libubox/jshn.sh
. /lib/functions.sh
user_exists admin || user_add admin
printf 'password\npassword\n' | /usr/bin/passwd admin
