# netgear/RAXE500/usr/bin/circlev2/shares/usr/bin/service/pingpongd/run
#!/bin/sh
exec $CIRCLE_ROOT/pingpongd -i $CIRCLE_BRIDGE
