# ubiquiti/U7-Pro/etc/uci-defaults/qca-mcsd
#!/bin/sh

VID=1
board=ap$(cat /tmp/sysinfo/board_name | awk -F 'ap' '{print $2}')
case "$board" in
ap-mp03.1*|ap-mi01*)
	VID=0
	;;
esac

uci -q batch << EOF
	set network.lan.multicast_querier=0
	set network.lan.igmp_snooping=0
	commit network
	set mcsd.config.SwitchLanVid=$VID
	uci commit mcsd
EOF

exit 0

