# netgear/RAX38/usr/lib/opkg/info/libcellwan.list
/opt/lantiq/etc/cell_wan_dev_intel_mbim.conf
/opt/lantiq/usr/sbin/wwancli
/opt/lantiq/www/cellularusimviewform.xml
/opt/lantiq/www/cellularview.xml
/etc/hotplug.d/usb/ltq-wwan-usb-lte-intel_dongles.hotplug
/opt/lantiq/xml/Cellular_data.xml
/opt/lantiq/www/cellularviewform.xml
/opt/lantiq/www/cellularusimview.xml
/etc/hotplug.d/usb/ltq-wwan-usb-lte-intel_at.hotplug
/opt/lantiq/etc/cell_wan_dev_intel_at.conf
/opt/lantiq/servd/lib/libcw.so
/opt/lantiq/www/cellularwan_menu.json
/opt/lantiq/usr/lib/libwwanfapi.so
/opt/lantiq/usr/lib/libcw.so
/usr/include/wwanfapimbim.h
/etc/hotplug.d/usb/ltq-wwan-usb-lte-intel_mbim.hotplug
/opt/lantiq/xml/Cellular_control.xml
/opt/lantiq/etc/init.d/intel_init_at.sh
/opt/lantiq/etc/init.d/intel_init_mbim.sh
/opt/lantiq/usr/sbin/dcd
/opt/lantiq/etc/init.d/intel_uninit.sh
/usr/include/wwanfapiat.h
