# netgear/RAX38/usr/lib/opkg/info/ltq-voip-common.list
/etc/init.d/ifx_load_voip_common_drivers.sh
/lib/modules/3.10.108/voip_timer_driver.ko
