# netgear/RAX38/usr/lib/opkg/info/ltq-evtlog.list
/etc/init.d/load_evtlog.sh
/opt/lantiq/bin/drv_ifxos.ko
/lib/modules/3.10.108/drv_event_logger.ko
/opt/lantiq/bin/drv_event_logger.ko
