# netgear/RAX38/usr/lib/opkg/info/pad.list
/etc/pad/plugins/2-ddr.plugin
/etc/pad/plugins/datapath.plugin
/etc/pad/plugins/23-pp_wifi_info.plugin
/opt/lantiq/usr/sbin/ugw_throughput.sh
/etc/pad/plugins/19-perf.plugin
/etc/pad/plugins/14-mcast.plugin
/etc/pad/plugins/17-ubi.plugin
/etc/pad/help/tune_help.txt
/etc/pad/plugins/22-pp_info.plugin
/etc/pad/plugins/25-dbfile.plugin
/ramdisk_copy/tmp/plugins/et/tools/pluginscript.sh
/etc/init.d/pad.sh
/opt/lantiq/usr/sbin/read_core_dump.sh
/ramdisk_copy/tmp/plugins/datapath.plugin
/ramdisk_copy/tmp/plugins/dp/tools/rmonclear.sh
/opt/lantiq/usr/sbin/cpu_usage.sh
/ramdisk_copy/tmp/plugins/et/tools/create-commands.sh
/etc/pad/plugins/4-dma.plugin
/opt/lantiq/usr/sbin/ltq_probe_storage.sh
/etc/pad/help/debug_help.txt
/etc/pad/plugins/8-coc.plugin
/etc/pad/plugins/20-pp_cli.plugin
/ramdisk_copy/tmp/plugins/dp/tools/create-commands.sh
/etc/pad/plugins/26-wave_info.plugin
/ramdisk_copy/tmp/plugins/endurance.plugin
/etc/pad/import_platform.plugin
/etc/pad/plugins/6-gphy.plugin
/ramdisk_copy/tmp/plugins/et/tools/endurance.sh
/etc/pad/plugins/11-qos.plugin
/etc/pad/help/diag_help.txt
/opt/lantiq/usr/sbin/qos_ethertype_to_pcp.sh
/etc/pad/plugins/29-wifi.plugin
/etc/pad/plugins/3-cbm.plugin
/opt/lantiq/usr/sbin/pad
/etc/pad/plugins/16-usb.plugin
/etc/pad/plugins/automation.plugin
/opt/lantiq/usr/sbin/qos_dmac_vlan_to_port.sh
/etc/pad/plugins/1-gswip.plugin
/etc/pad/plugins/24-ipv6_info.plugin
/opt/lantiq/usr/sbin/pad_cleanup.sh
/etc/pad/plugins/datapath_debug.plugin
/etc/pad/plugins/7-system.plugin
/ramdisk_copy/tmp/plugins/dp/tools/pluginscript.sh
/etc/pad/plugins/27-crash_dump.plugin
/etc/pad/plugins/5-ppa.plugin
/opt/lantiq/usr/sbin/ltq_probe_tftpserver.sh
/etc/pad/plugins/10-networking.plugin
/opt/lantiq/usr/sbin/pp_cli_cmd.sh
/etc/pad/plugins/9-flash.plugin
/etc/pad/help/datapath_help.txt
/etc/pad/help/endurance_help.txt
/etc/pad/plugins/12-vmb.plugin
/ramdisk_copy/tmp/plugins/dp/tools/datapath.sh
/opt/lantiq/usr/sbin/qos_dscp_tc_to_pcp.sh
/etc/pad/plugins/21-ugw_throughput.plugin
/opt/lantiq/usr/sbin/platform_discover.sh
/etc/pad/plugins/endurance.plugin
/etc/pad/plugins/15-framework.plugin
