# netgear/RAX38/usr/lib/opkg/info/polld.list
/etc/init.d/polld
/opt/lantiq/bin/polld
