# netgear/RAX38/usr/lib/opkg/info/pudil.list
/usr/bin/pudil
/etc/init.d/dil_monitor
/etc/init.d/dil_pega.init
/usr/bin/pudild
/usr/lib/libpudil.so
/usr/sbin/pudil_update_devStatusByClid.sh
/etc/init.d/pucb.sh
/usr/bin/pucb_ac
/usr/bin/pucb_dil
/usr/bin/pucb_conndev
