# netgear/RAX38/usr/lib/opkg/info/putmd.list
/usr/sbin/putm
/usr/sbin/putmd
/etc/init.d/putmd
