# netgear/RAX38/usr/lib/opkg/info/session-mgmt.list
/opt/lantiq/etc/rc.d/sess_mgmt.sh
/opt/lantiq/usr/sbin/ppasessmgmt
/etc/init.d/sess_mgmt
/opt/lantiq/usr/sbin/ppasessmgmtd
/opt/lantiq/etc/rc.d/flow-aware-qos.sh
/opt/lantiq/etc/sess_mgmt.conf
