# netgear/RS700S/usr/bin/circlev2/shares/usr/bin/service/ipaddr_manager/run
#!/bin/sh
exec $CIRCLE_ROOT/ipaddr_manager -d
