# tplink_en/TL-WR941HP/etc/hotplug.d/button/52-wps-supplicant
#
# Copyright (c) 2014, The Linux Foundation. All rights reserved.
#
#  Permission to use, copy, modify, and/or distribute this software for any
#  purpose with or without fee is hereby granted, provided that the above
#  copyright notice and this permission notice appear in all copies.
#
#  THE SOFTWARE IS PROVIDED "AS IS" AND THE AUTHOR DISCLAIMS ALL WARRANTIES
#  WITH REGARD TO THIS SOFTWARE INCLUDING ALL IMPLIED WARRANTIES OF
#  MERCHANTABILITY AND FITNESS. IN NO EVENT SHALL THE AUTHOR BE LIABLE FOR
#  ANY SPECIAL, DIRECT, INDIRECT, OR CONSEQUENTIAL DAMAGES OR ANY DAMAGES
#  WHATSOEVER RESULTING FROM LOSS OF USE, DATA OR PROFITS, WHETHER IN AN
#  ACTION OF CONTRACT, NEGLIGENCE OR OTHER TORTIOUS ACTION, ARISING OUT OF
#  OR IN CONNECTION WITH THE USE OR PERFORMANCE OF THIS SOFTWARE.
#

WIFI_RE_BUTTON_PID="/var/run/wifi_re_button.pid"

if [ "$ACTION" = "released" -a "$BUTTON" = "re" ]; then
        if [ "$(uci get factory.factorymode.enable)" = "yes" ]; then
            touch /tmp/button_re_check
        fi
        if [ "$SEEN" -lt 2 ]; then
                [ -e ${WIFI_RE_BUTTON_PID} ] && {
                        local pid=$(cat ${WIFI_RE_BUTTON_PID})
                        echo "$ACTION:$BUTTON:$SEEN: Kill pid $pid." > /dev/console
                        [ -n "$pid" ] && kill -9 $pid
                        rm -f ${WIFI_RE_BUTTON_PID}
                }
        fi
fi

if [ "$ACTION" = "pressed" -a "$BUTTON" = "re" ]; then
        if [ "$(uci get factory.factorymode.enable)" = "yes" ]; then
            return
        fi
        /sbin/wifi_re_button &
fi
