# tplink_en/Archer AX11000/usr/lib/opkg/info/product-file.list
/sbin/board_restore
/lib/preinit/06_init_insmod_usbhost_scsi
/lib/wifi/brcm_wlan_config.txt
/usr/sbin/wps_status
/sbin/board_factory
/usr/lib/.gitignore
/usr/sbin/status
/etc/hotplug.d/net/10-wl_cpu_affinity
/lib/wifi/brcmwifi.sh
/etc/init.d/wireless
/sbin/wifi_button
/lib/print_server/core_print_server.sh
/etc/hotplug.d/net/21-switch_ports
/lib/wifi/wps_button
/lib/led_switch
/sbin/wifi_firm
/etc/hotplug.d/net/20-client_mgmt
/etc/hotplug.d/usb/21-usb_reset
/lib/preinit/05_init_interfaces_brcm
/lib/preinit/71_init_etc_ramfs
/lib/print_server/printer_driver.sh
/usr/sbin/lan_connected
/etc/hotplug.d/block/60-enable_gro
/etc/hotplug.d/net/01-wds
/lib/nat/nat_hw.sh
/etc/hotplug.d/net/22-add_neigh
/usr/sbin/wan_connected
/usr/lib/lua/luci/model/dstRuleTbl.lua
/etc/hotplug.d/iface/02-wireless-vlan
/etc/hotplug.d/usb/20-usb_cpu_affinity
