# tplink_en/Archer AX11000/usr/lib/opkg/info/speed-test.list
/etc/speed-test/settings.txt
/usr/bin/speed-test
/etc/init.d/speedtest
/etc/init.d/speedtest_status
/usr/sbin/speedtest_polling
/usr/sbin/speedtest
/etc/speed-test/get_service_status.sh
/etc/speed-test/speed_test.sh
