# tplink_en/Archer AX11000/etc/init.d/postcenter
#!/bin/sh /etc/rc.common

START=50

start()
{
        service_start /usr/sbin/postcenter &
}

stop()
{
        service_stop /usr/sbin/postcenter
}

