# tplink_en/Archer AX6000/usr/lib/opkg/info/bluetooth.list
/usr/lib/lua/bluetooth.lua
/etc/init.d/bluetooth
/usr/sbin/test_ble
/usr/sbin/bluetooth
/etc/init.d/bsa_server
/etc/config/bluetooth
