# tplink_en/Archer AX6000/usr/lib/opkg/info/ledcenter.list
/usr/bin/ledcli
/usr/bin/ledcenter
/etc/config/ledcenter
/etc/init.d/ledcenter
