# tplink_en/Touch P5/etc/inittab
::sysinit:/etc/rc.d/rcS >/dev/console 2>&1
::respawn:/sbin/getty ttyS0 115200
::shutdown:/bin/umount -a
