# netgear/RAX38/etc/init.d/dal_pega.preinit
#!/bin/sh /etc/rc.common
START=16

start() {
	mkdir -p /etc/d2
	mkdir -p /etc/dal
	ln -s /etc/dal /tmp/arlo
	ln -s /etc/dal /tmp/dal
}

