# netgear/RAX35/usr/lib/opkg/info/base-files.list
/etc/cmv_batch
/etc/openwrt_version
/etc/tr69
/opt/lantiq/etc/switchports.conf
/ramdisk_copy/etc/inetd.conf
/etc/banner.failsafe
/ramdisk_copy/tmp/u_lock
/etc/init.d/disable_bridge_acceleration.sh
/etc/fastmode/runwifi.sh
/etc/rc.d/udhcpc_ad
/bin/config_generate
/var
/etc/nsswitch.conf
/etc/opkg/keys/de98a2dd1d0f8a07
/etc/preinit
/etc/rc.d/init.d
/etc/host.conf
/etc/init.d/init_verify
/lib/upgrade/platform.sh
/usr/sbin/version.sh
/etc/shells
/lib/preinit/10_sysinfo
/lib/preinit/70_initramfs_test
/etc/device_info
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/lib/firmware/ppe_fw.bin
/lib/preinit/10_indicate_preinit
/ramdisk_copy/etc/TZ
/etc/profile.d/ugw_framework_env.sh
/etc/modules-production.d/43-dc_mode0-xrx500
/usr/sbin/mini_fo.sh
/etc/ripd.conf
/etc/dnrd
/lib/functions/uci-defaults-new.sh
/etc/udhcpd.conf
/etc/config/system
/mnt/data/fs
/etc/hotplug.d/usb/ltq-usb-led.hotplug
/etc/ipsec.conf
/etc/fstab
/etc/sysctl.conf
/sbin/sysupgrade
/etc/fastmode/prepare.sh
/lib/preinit/99_10_run_init
/etc/device_table.txt
/tmp
/etc/init.d/radvd
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/etc/shadow
/etc/modules-production.d/42-directconnect_datapath
/etc/rc.common
/lib/preinit/10_indicate_failsafe
/etc/rc.d/wlan_discover
/init
/ramdisk_copy/etc/dsl_api/cmv_scripts
/etc/spmd.pwd
/lib/preinit/42_athfix
/etc/init.d/passwd.sh
/etc/ipsec.secrets
/usr/sbin/wlan_hw_revision.sh
/ramdisk_copy/etc/ppp/resolv.conf
/etc/init.d/inetd
/lib/lantiq.sh
/etc/modules.conf
/lib/preinit/80_mount_root
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/rc.d/rebootcpe.sh
/etc/rc.d/invoke_upgrade.sh
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/usr/sbin/overlayfs.sh
/etc/snmp
/firmware
/etc/update_calibration.sh
/ramdisk_copy/etc/ilmid/ilmid.conf
/etc/rc.button/ledswitch
/etc/rc.d/config.sh
/etc/rc.button/wps
/etc/rc.button/failsafe
/usr/sbin/factorycfg.sh
/sbin/modprobe.sh
/etc/rc.d/killproc
/etc/init.d/set_system_config.sh
/etc/services
/lib/functions/system.sh
/etc/rc.d/checkpid
/bin/board_detect
/lib/upgrade/common.sh
/etc/rc.button/power
/etc/sysupgrade.conf
/flash
/etc/hotplug.d/firmware/10-rt2x00-eeprom
/etc/fastmode/restore.sh
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/ramdisk_copy/etc/services
/etc/config/system_RAX35
/etc/init.d/version.sh
/lib/preinit/40_run_failsafe_hook
/etc/group
/etc/rc.button/reset
/lib/preinit/30_failsafe_wait
/etc/hotplug.d/button/10-generic.sh
/etc/inetd.conf
/sbin/firstboot
/etc/rc.d/udhcpc.script
/bin/login.sh
/ramdisk_copy/tmp/s_lock
/etc/zebra.conf
/usr/sbin/getphymem
/etc/hotplug.d/net/00-sysctl
/etc/hosts
/lib/functions.sh
/etc/uci-defaults/11_migrate-sysctl
/lib/functions/leds.sh
/sbin/clishshell
/usr/sbin/ugw_debug_print.sh
/rom/note
/etc/inittab
/etc/openwrt_release
/etc/rc.d/disable_RGMII_SSC.sh
/sbin/wifi
/etc/rc.conf
/etc/ilmid
/etc/profile
/etc/profile.arch
/lib/firmware/09a9/aca_fw.bin
/lib/preinit/99_10_failsafe_login
/etc/fastmode/set_system_config.sh
/lib/functions/service.sh
/lib/upgrade/keep.d/base-files-essential
/sbin/hotplug
/lib/preinit/02_default_set_state
/usr/sbin/vol_mgmt
/etc/mtab
/etc/init.d/boot
/etc/init.d/sys_opt.sh
/ramdisk_copy/tmp/system_status
/etc/init.d/usb_host_device.sh
/ramdisk_copy/etc/ripd.conf
/ramdisk_copy/etc/hosts
/etc/extmark.sh
/lib/functions/preinit.sh
/etc/modules-production.d/45-swa_stack_al
/bin/ipcalc.sh
/etc/fastmode/enable.sh
/mnt/usb
/etc/diag.sh
/etc/hostapd.conf
/etc/rc.d/udhcpc_lan.script
/etc/init.d/udhcpc
/etc/codeclist
/etc/rc.button/wifi
