# netgear/RAX35/usr/lib/opkg/info/csd.list
/opt/lantiq/usr/sbin/capitest
/usr/sbin/csdutil
/opt/lantiq/usr/sbin/csd
/etc/init.d/csd
/opt/lantiq/usr/lib/libcapi.so
