# netgear/RAX35/usr/lib/opkg/info/dil.list
/etc/init.d/dil
/usr/lib/libdalcjson.so
/usr/lib/libd2dil.so
/etc/dil/file1.txt
/usr/bin/diltest
/etc/dil/file2.txt
/etc/init.d/dilboot
/usr/bin/fing_dil
/usr/bin/jsontest
