# netgear/RAX35/usr/lib/opkg/info/ltq-ppmd.list
/opt/lantiq/bin/ppmd
/etc/init.d/ppmd.sh
