# netgear/RAX35/usr/lib/opkg/info/ltq-temp.list
/lib/modules/3.10.108/ltq_temp.ko
/etc/init.d/ltq_temp.sh
