# netgear/RAX35/usr/lib/opkg/info/ltq-pm-util.list
/opt/lantiq/bin/pm_util
/etc/profile.d/ltq_pm.sh
/etc/init.d/reset_regulator_sw_control.sh
/etc/init.d/pm_util.sh
