# netgear/RAX35/usr/lib/opkg/info/ntgrutichk.list
/etc/init.d/ntgrutichk
/usr/bin/ntgrutichk
