# netgear/RAX35/usr/lib/opkg/info/pu_attachDevd.list
/usr/bin/pu_attachDevd
/etc/init.d/pu_attachDevd
