# netgear/RAX35/usr/lib/opkg/info/pu_sshEnabled.list
/usr/bin/pu_sshEnabled
/etc/init.d/pu_sshEnabled
