# netgear/RAX35/usr/lib/opkg/info/pu_telnetEnabled.list
/etc/init.d/pu_telnetEnabled
/usr/bin/pu_telnetEnabled
