# netgear/RAX35/usr/lib/opkg/info/pucgi.list
/opt/vendor/www2/cgi/traffic-control
/opt/vendor/www2/cgi/pucgid
/opt/vendor/www2/vendorwhitehost.conf
/opt/vendor/www2/www2-tnc.tgz
/opt/vendor/www2/cgi/pucgi
/etc/init.d/openvpn1
/etc/init.d/openvpn2
/opt/vendor/www2/lighttpd4.conf
/opt/vendor/www2/www2.tgz
/opt/vendor/www2/lighttpd5.conf
/etc/init.d/pucgid
/opt/vendor/www2/www2-lang.tgz
/etc/init.d/lighttpd4
/etc/init.d/lighttpd5
