# netgear/RAX35/usr/lib/opkg/info/puhttpsniff.list
/etc/init.d/puhttpsniff
/opt/lantiq/bin/puhttpsniff
