# netgear/RAX35/usr/lib/opkg/info/xagent.list
/etc/xagent/xagent/xagent.prod.env
/opt/xagent/xagent_control
/etc/xagent/remGenie/remGenie.qa.env
/etc/xagent/xagent/switch.sh
/etc/xagent/upagent/upagent.beta.env
/etc/xagent/upagent/upagent.qa2.env
/opt/rcagent/scripts/comm.sh
/www/cgi-bin/genie.cgi
/usr/lib/libreadycloud.so
/opt/leafp2p/leafp2p.tar.gz
/etc/xagent/upagent/upagent.prod.env
/www/cgi-bin/readycloud_control.cgi
/etc/xagent/remGenie/remGenie.qa2.env
/etc/xagent/remGenie/switch.sh
/opt/rcagent/scripts/alias.sh
/etc/xagent/bst/bst.qa.env
/opt/remote/run_remote
/etc/xagent/upagent/upagent.ftrial.env
/opt/rcagent/scripts/env_nvram.sh
/etc/xagent/upagent/upagent.staging.env
/etc/xagent/xcenvsw.sh
/etc/xagent/upagent/upagent.qa.env
/opt/leafp2p/leafp2p
/etc/xagent/remGenie/remGenie.staging.env
/opt/rcagent/scripts/unregister.sh
/etc/xagent/xagent/xagent.dev.env
/etc/xagent/xagent/xagent.ftrial.env
/opt/rcagent/cgi_processor
/opt/remote/remote
/opt/rcagent/scripts/register.sh
/opt/xagent/genie_handler
/opt/xagent/xagent
/etc/xagent/upagent/switch.sh
/etc/preinit.d/xagent_preinit.sh
/opt/rcagent/run_server.sh
/etc/xagent/bst/bst.prod.env
/opt/xagent/certs/.ca-bundle-mega.crt.swp
/etc/xagent/remGenie/remGenie.beta.env
/etc/xagent/xagent/xagent.staging.env
/etc/xagent/bst/bst.staging.env
/usr/bin/evsinkd2
/opt/rcagent/rcagentd
/etc/xagent/xagent/xagent.qa2.env
/opt/rcagent/downloader
/etc/xagent/upagent/upagent.dev.env
/etc/xagent/remGenie/remGenie.prod.env
/etc/xagent/bst/switch.sh
/opt/rcagent/scripts/env.sh
/etc/xagent/xagent/xagent.qa.env
/opt/rcagent/nas_service
/etc/xagent/remGenie/remGenie.ftrial.env
/etc/xagent/xagent/xagent.beta.env
/opt/xagent/run-xagent.sh
/opt/xagent/certs/ca-bundle-mega.crt
/etc/xagent/remGenie/remGenie.dev.env
/etc/xagent/bst/bst.dev.env
/opt/rcagent/cgi/rccommand.cgi
/etc/init.d/XAgent
