# netgear/RAX40/usr/lib/opkg/info/ltq-cgroups.list
/usr/bin/cgroup_daemon
/usr/bin/cgroup_notifier
/etc/preinit.d/ltq_cgroups_preinit.sh
/etc/init.d/ltq_cgroups_startup.sh
/etc/cgroups.conf
