# netgear/RAX40/usr/lib/opkg/info/minissdpd.list
/etc/init.d/minissdpd
/usr/bin/minissdpd
/etc/config/minissdpd
