# netgear/RAX40/usr/lib/opkg/info/ntgr_wifischedule.list
/usr/bin/ntgr_wifischedule
/etc/init.d/ntgr_wifischedule
