# netgear/RAX40/usr/lib/opkg/info/ntgr_ra_dce.list
/etc/init.d/ntgr_ra_dce
/usr/sbin/NTGR_RA_memCalcUsageHigh.sh
/usr/sbin/NTGR_RA_cpuUsage.sh
/usr/bin/ntgr_ra_dce
/usr/sbin/NTGR_RA_memUsage.sh
/usr/sbin/NTGR_RA_disconnectCount.sh
/etc/init.d/ntgr_bootcheck
/usr/sbin/NTGR_RA_cpuCalcAllAve.sh
