# netgear/RAX40/usr/lib/opkg/info/pudal.list
/etc/init.d/dal_pega.init
/etc/init.d/dal_pega.preinit
/usr/bin/pudal
/etc/init.d/~dal_boot_notify
