# netgear/RAX40/usr/lib/opkg/info/samba36-server.list
/etc/config/samba
/ramdisk_copy/flash/samba/smb.conf
/usr/sbin/nmbd
/etc/samba/smbpasswd
/etc/samba/valid.dat
/etc/init.d/samba
/etc/samba/upcase.dat
/usr/sbin/smbpasswd
/usr/sbin/smbd
/etc/functions.sh
/lib/upgrade/keep.d/samba36-server
/usr/sbin/samba_multicall
/etc/smb.conf
/etc/samba/lowcase.dat
