# netgear/RAX40/usr/lib/opkg/info/servd.list
/opt/lantiq/usr/sbin/servd
/opt/lantiq/usr/sbin/slloglevel
/etc/init.d/servd
