# netgear/RAX40/usr/lib/opkg/info/soapservd.list
/etc/init.d/soapservd2_for_xagent
/etc/init.d/soapservd
/usr/bin/soapservd
