# netgear/RAX40/usr/lib/opkg/info/syslog-ng.list
/etc/init.d/syslog-ng
/opt/lantiq/etc/syslog-ng.conf
/usr/sbin/syslog-ng
/lib/upgrade/keep.d/syslog-ng
