# ubiquiti/UAP-AC-M-PRO/usr/lib/opkg/info/ubnt-lede-files.list
/usr/etc/rc.d/rc.e515
/usr/etc/default-e302.cfg
/usr/etc/default_qca933x.cfg
/usr/etc/default-e507.cfg
/usr/etc/default.cfg
/usr/etc/rc.d/rc.disable_wb_cache
/usr/etc/default-nossid-e582.cfg
/usr/etc/default-nossid-e547.cfg
/bin/walled_action.sh
/bin/syslogd_wrapper.sh
/usr/etc/default-ec20.cfg
/usr/etc/rc.d/rc.sysinit
/usr/etc/common.sh
/usr/lib/version
/usr/etc/default-e517.cfg
/usr/etc/default_nossid_qca956x.cfg
/usr/etc/default-nossid-e592.cfg
/usr/etc/default-nossid-e557.cfg
/usr/etc/udhcpc/udhcpc_ip_only
/usr/etc/.build
/usr/etc/default-e587.cfg
/usr/etc/default-nossid-e532.cfg
/usr/etc/default-dca8.cfg
/usr/etc/passwd
/usr/etc/default_nossid_ar724x.cfg
/usr/etc/ethertypes
/usr/lib/semver
/usr/etc/hotplug2.rules
/usr/etc/default-e562.cfg
/usr/etc/default-e527.cfg
/bin/syswrapper.sh
/usr/etc/default_nossid_qca933x.cfg
/usr/etc/default_qca9342.cfg
/usr/etc/default-nossid-e515.cfg
/usr/etc/platdep_funcs.sh
/usr/etc/host.conf
/usr/etc/default_nossid_ar9344.cfg
/usr/etc/qca-full-offload
/usr/etc/default-e502.cfg
/etc/startup.list
/usr/etc/default-nossid-e567.cfg
/lib/preinit/11_config_tcp_timestamps
/usr/etc/protocols
/usr/etc/fstab
/usr/etc/rc.d/rc.modules
/usr/etc/udhcpc_services
/usr/etc/rc.d/rc.plugins
/usr/etc/default-nossid-e302.cfg
/usr/etc/default_ar724x.cfg
/usr/etc/default-nossid-e507.cfg
/bin/show_node
/usr/etc/login.defs
/usr/etc/inittab
/usr/etc/default-e572.cfg
/usr/etc/default-e537.cfg
/usr/etc/startup.list
/usr/etc/default_ar9344.cfg
/usr/etc/profile
/usr/etc/services
/usr/etc/shells
/usr/etc/default-e512.cfg
/bin/htb
/lib/preinit/10_mmcblk
/usr/etc/syswrapper.sh
/bin/show_nt
/etc/init.d/ubnt_sysinit
/lib/preinit/99_20_ubnt_disable_services
/etc/uap-sysctl.conf
/usr/etc/activate_tunnel.sh
/usr/etc/default-nossid-e517.cfg
/usr/etc/default_nossid.cfg
/usr/etc/default-e582.cfg
/usr/etc/default-e547.cfg
/usr/etc/udhcpc/udhcpc
/usr/etc/apscan.sh
/usr/etc/rc.d/rc
/usr/etc/.version
/usr/etc/default_nossid_qca9342.cfg
/bin/ubnt-rssimon
/lib/preinit/99_21_ubnt_ubntconf
/bin/radartoolw
/usr/etc/default-nossid-e587.cfg
/usr/etc/rc.d/rc.identity
/usr/etc/default-nossid-dca8.cfg
/usr/etc/qca-package
/usr/etc/banner
/usr/etc/default_single.cfg
/usr/etc/default-nossid-e562.cfg
/usr/etc/default-nossid-e527.cfg
/usr/etc/group
/usr/etc/default-e592.cfg
/usr/etc/default-e557.cfg
/usr/etc/default-nossid-e502.cfg
/bin/mmc_info
/usr/lib/build.properties
/usr/etc/default-e532.cfg
/usr/etc/hosts
/usr/etc/default_qca956x.cfg
/usr/etc/rc.d/rc.platform
/etc/init.d/ubnt_plugins
/bin/pktgen.sh
/usr/etc/cfg80211_support
/usr/etc/default-e515.cfg
/usr/etc/default-nossid-e572.cfg
/usr/etc/default-nossid-e537.cfg
/usr/etc/default-e567.cfg
/usr/etc/default-nossid-e512.cfg
/usr/etc/led_locate.sh
