# tplink_en/Archer GX90/usr/lib/opkg/info/lp5523ctrl.list
/etc/tsched.d/leds_schedule_allday
/etc/tsched.d/leds_schedule
/etc/sbin/leds_night_mode
/etc/init.d/lp5523ctrl
/usr/bin/lp5523ctrl
