# tplink_en/Archer GX90/usr/lib/opkg/info/tm-shn.list
/usr/bin/tm-shn/moni.sh
/usr/bin/tm-shn/tdts.ko
/usr/bin/tm-shn/rule_schema.trf
/usr/bin/tm-shn/list_security_rules.db
/usr/bin/tm-shn/rule.trf
/usr/bin/tm-shn/iqos-setup.sh
/usr/bin/tm-shn/sample.bin
/etc/init.d/tm_shn_sec
/usr/bin/tm-shn/libshn_pctrl.so
/lib/libshn_pctrl.so
/usr/bin/tm-shn/lic-setup.sh
/usr/bin/tm-shn/meta_en-US.dat
/usr/bin/tm-shn/pre_lic.key
/usr/bin/tm-shn/libtdts_ctrl.so
/lib/libtdts_ctrl.so
/usr/bin/tm-shn/wbl.conf
/usr/bin/tm-shn/wred
/etc/init.d/tm_shn_pc
/usr/bin/tm-shn/liblicop.so
/etc/init.d/tm_shn_qos
/usr/bin/tm-shn/clean-cache.sh
/usr/bin/tm-shn/dc_monitor.sh
/usr/bin/tm-shn/wred-setup.sh
/usr/bin/tm-shn/wred_set_wbl
/usr/bin/tm-shn/setup.sh
/usr/bin/tm-shn/shn_ctrl
/usr/bin/tm-shn/tdts_udb.ko
/usr/bin/tm-shn/patrol_tq.conf
/lib/liblicop.so
/usr/bin/tm-shn/qos.conf
/usr/bin/tm-shn/feedbacktest.cert
/usr/bin/tm-shn/libshn_utils.so
/usr/bin/tm-shn/app_patrol.conf
/lib/libshn_utils.so
/usr/bin/tm-shn/libopenlib.so
/usr/bin/tm-shn/list_pc_filter_apps.db
/usr/bin/tm-shn/tdts_ctrl
/usr/bin/tm-shn/wred.conf
/lib/libopenlib.so
/usr/bin/tm-shn/dcd
/usr/bin/tm-shn/tcd
/usr/bin/tm-shn/list_security_wrs_cats.db
/usr/bin/tm-shn/tdts_udbfw.ko
/etc/init.d/tm_shn
/usr/bin/tm-shn/wred_set_conf
/usr/bin/tm_shn
/usr/bin/tm-shn/tcd_monitor.sh
/usr/bin/tm-shn/timestamp.20200805
/usr/bin/tm-shn/gen_lic
/usr/bin/tm-shn/sig_version
