# netgear/RAX5/usr/lib/opkg/info/bst.list
/etc/init.d/bst
/usr/bin/bsttest
/usr/bin/bst_daemon
/usr/lib/libd2bst.so
/usr/lib/libd2bst.a
