# netgear/RAX5/usr/lib/opkg/info/lan-monitor.list
/etc/config/vpndev
/sbin/hostname-monitor
/sbin/lan-monitor
/sbin/vpn-monitor
/etc/init.d/lan-monitor
/usr/sbin/neigh.sh
