# netgear/RAX5/usr/lib/opkg/info/led_control.list
/usr/bin/led_control
/etc/init.d/led
/etc/hotplug.d/phy/00-ledstate
