# netgear/RAX5/usr/lib/opkg/info/ntgr_traffic_meterd.list
/etc/traffic_meter/dnsmasq_tmRedirecHttpHttps.conf
/etc/init.d/traffic_meterd
/usr/sbin/tm_reload.sh
/usr/sbin/ntgr_trafficmeter
