# netgear/RAX5/usr/lib/opkg/info/ntgr_ra_iot.list
/etc/rabin/enable_ra_log
/usr/bin/puraUpdate
/usr/bin/NTGR_RA_cpuUsage.sh
/etc/rabin/ntgr_ra_iot
/usr/bin/ra_cli
/usr/bin/ra_nvram
/usr/bin/NTGR_RA_cpuCalcAllAve.sh
/etc/rabin/raePolicy.json
/etc/init.d/ntgr_ra_iot
/lib/pega/rae_kernel_monitor.ko
