# netgear/RAX5/usr/lib/opkg/info/pot.list
/etc/init.d/potd
/usr/sbin/pot
/usr/sbin/createpotfile.sh
