# netgear/RAX5/usr/lib/opkg/info/pushntf.list
/usr/bin/dalh
/usr/bin/pushntftest
/etc/init.d/DALH
