# netgear/RAX5/usr/lib/opkg/info/rax5.list
/etc/config/omcproxy
/lib/netifd/proto/dhcp.sh
/etc/debug.conf
/etc/init.d/odhcpd
/etc/config/WlanGuestSchedule
/etc/config/WebServiceManagement
/etc/openvpncfg/ovpncfg-stop.sh
/lib/netifd/proto/6to4.sh
/etc/openvpncfg/ovpncfg.sh
/etc/config/netgear
/lib/ramips.sh
/etc/firewall.d/urlfilter
/lib/preinit/05_set_iface_mac_mediatek
/etc/init.d/network
/etc/firewall.d/dnshijack
/etc/config/loginpwd
/etc/config/blocksch
/lib/netifd/proto/v6passthru.sh
/etc/firewall.d/nat_cone
/lib/upgrade/platform.sh
/usr/sbin/ssh_telnet_Enabled.sh
/etc/init.d/qos
/usr/sbin/installEvent_wifi_status.sh
/etc/pega-defaults/run_script-2
/etc/pega-defaults/run_script-4
/etc/pega-defaults/run_script-9
/lib/wifi/mt7615e.lua
/etc/init.d/vlan
/usr/sbin/quagga.init
/etc/openvpncfg/server_tun.conf.tmpl
/lib/netifd/proto/l2tp.sh
/etc/init.d/cron
/lib/functions/pega.sh
/etc/hotplug.d/phy/01-wanstate
/usr/sbin/inet_connect.sh
/etc/init.d/email
/etc/init.d/wifi
/usr/sbin/check_ipconflict.sh
/etc/config/upnpd
/usr/sbin/ookla
/etc/uci-defaults/32_gen_NTGR_WiFi
/usr/lib/lua/shuci.lua
/lib/netifd/dhcp.script
/usr/sbin/debug.sh
/etc/config/system
/etc/init.d/firewall
/etc/config/rpcd
/etc/init.d/blockservice
/etc/openvpncfg/genvpnclientlist.sh
/etc/sysctl.conf
/sbin/sysupgrade
/etc/config/ipLocation
/lib/netifd/proto/dhcpv6.sh
/etc/config/ucitrack
/usr/lib/ddns/dynamic_dns_updater.sh
/etc/openvpncfg/client_tap.ovpn.tmpl
/etc/wireless/mediatek/mt7915.dbdc.b0.dat
/etc/firewall.d/vpn_devname
/etc/config/openvpn
/etc/config/vlan
/etc/init.d/pega_pre
/usr/sbin/wifi_guest_schedule_2g.sh
/etc/init.d/dropbear
/etc/pega-defaults/patch_db-1
/etc/hotplug.d/phy/02-wan_v6state
/etc/shadow
/etc/pega-defaults/patch_db-3
/etc/pega-defaults/patch_db-4
/etc/pega-defaults/patch_db-5
/etc/pega-defaults/patch_db-6
/etc/pega-defaults/patch_db-7
/etc/pega-defaults/patch_db-8
/etc/init.d/ddns
/etc/config/blocksite
/etc/hotplug.d/iface/95-ddns
/etc/board.d/02_network
/etc/firewall.d/blockSrv
/etc/init.d/wan_up
/usr/sbin/inet_disconn.sh
/usr/sbin/lighttp_restart.sh
/usr/sbin/email_log.sh
/etc/config/routing
/usr/sbin/debug_info.sh
/etc/init.d/pega_post
/etc/config/firewall
/lib/preinit/80_mount_root
/etc/hotplug.d/phy/00-phystate
/etc/config/dhcp
/etc/wireless/mediatek/mt7915.dbdc.b1.dat
/etc/hotplug.d/iface/50-miniupnpd
/etc/firewall.d/routerlogin
/etc/pega-defaults/patch_db-10
/etc/uci-defaults/30_uboot-envtools
/etc/fw_env.config
/etc/config/dropbear
/etc/config/ddns
/lib/upgrade/fwtool.sh
/etc/firewall.user
/usr/lib/ddns/dynamic_dns_functions.sh
/etc/uci-defaults/31_sync-pdata
/usr/sbin/inet_check.sh
/etc/openvpncfg/client_tun.ovpn.tmpl
/usr/sbin/blockschedule_on.sh
/usr/sbin/factory_reset.sh
/etc/hotplug.d/iface/90-auto-time-sync
/lib/upgrade/common.sh
/etc/openvpncfg/ovpnlog.sh
/usr/sbin/blockschedule_off.sh
/usr/sbin/debug_output.sh
/etc/init.d/openvpn
/etc/config/landev
/etc/config/tm
/etc/ssl.conf
/lib/network/switch.sh
/etc/rc.button/reset
/etc/init.d/dnsmasq
/etc/uci-defaults/netgear.default
/etc/init.d/blocksite
/usr/sbin/wifi_guest_schedule_5g.sh
/usr/sbin/lan_linkdownup.sh
/etc/init.d/miniupnpd
/etc/ppp/options
/usr/sbin/ntgr_ipLocation.sh
/lib/config/uci.sh
/etc/config/network
/etc/config/logs
/usr/sbin/sysntpd_restart.sh
/usr/sbin/wifischedule.sh
/etc/hotplug.d/iface/99-checkip
/etc/firewall.d/firewall.upnpd
/etc/init.d/quagga
/usr/sbin/wan6_autodetect.sh
/etc/config/qos
/etc/openvpncfg/updateClientCfg.sh
/etc/config/armor
/etc/hotplug.d/ntp/30-ntp-pot
/etc/init.d/omcproxy
/usr/sbin/delay_freset.sh
/etc/firewall.d/accessCtrl
/sbin/mtk_factory_rw.sh
/lib/netifd/dhcpv6.script
/lib/netifd/proto/ppp.sh
/etc/firewall.d/dosProtection
/usr/sbin/genkey.sh
/usr/lib/ddns/update_no-ip_com.sh
/etc/config/email
/etc/passwd
/etc/init.d/access_ctrl
/etc/hotplug.d/wifi/00-bridge_state
/etc/openvpncfg/server_tap.conf.tmpl
/etc/pega-defaults/run_script-10
/etc/init.d/boot
/etc/hotplug.d/iface/20-inetstatus
/etc/init.d/blocksch
/usr/sbin/ddns_config.sh
/etc/firewall.d/dmz
/usr/sbin/wifi_schedule_radio.sh
/etc/dropbear/authorized_keys
/usr/sbin/soap_ntgr_speedtest.sh
