# netgear/RAX5/usr/lib/opkg/info/remGenie.list
/www/cgi-bin/genie.cgi
/etc/init.d/remGenie
/usr/include/libusi.h
/usr/lib/libusi.so
