# netgear/RAX5/usr/lib/opkg/info/rpcd.list
/etc/init.d/rpcd
/sbin/rpcd
/usr/share/rpcd/acl.d/unauthenticated.json
