# netgear/RAX5/usr/lib/opkg/info/xagent.list
/etc/xagent/xagent/xagent.prod.env
/etc/xagent/remGenie/remGenie.qa.env
/opt/xagent/xagent_control
/etc/xagent/xagent/switch.sh
/etc/xagent/upagent/upagent.beta.env
/etc/xagent/upagent/upagent.qa2.env
/opt/etc/usi/usi.profile
/etc/xagent/csh/csh.qa.env
/etc/xagent/upagent/upagent.prod.env
/etc/xagent/csh/csh.staging.env
/etc/xagent/csh/csh.prod.env
/etc/xagent/csh/csh.dev.env
/etc/xagent/remGenie/remGenie.qa2.env
/etc/xagent/remGenie/switch.sh
/opt/usi/usi
/etc/xagent/bst/bst.qa.env
/etc/xagent/upagent/upagent.ftrial.env
/etc/xagent/upagent/upagent.staging.env
/opt/etc/usi.conf
/opt/etc/usi/profiles/d2.profile
/etc/xagent/xcenvsw.sh
/etc/xagent/upagent/upagent.qa.env
/etc/xagent/remGenie/remGenie.staging.env
/etc/xagent/csh/switch.sh
/etc/xagent/xagent/xagent.dev.env
/etc/xagent/xagent/xagent.ftrial.env
/opt/xagent/genie_handler
/opt/xagent/xagent
/etc/xagent/upagent/switch.sh
/etc/xagent/bst/bst.prod.env
/etc/xagent/remGenie/remGenie.beta.env
/etc/xagent/xagent/xagent.staging.env
/etc/xagent/bst/bst.staging.env
/usr/bin/evsinkd2
/etc/xagent/xagent/xagent.qa2.env
/etc/xagent/upagent/upagent.dev.env
/etc/xagent/remGenie/remGenie.prod.env
/etc/xagent/bst/switch.sh
/etc/xagent/xagent/xagent.qa.env
/etc/xagent/remGenie/remGenie.ftrial.env
/etc/xagent/xagent/xagent.beta.env
/opt/xagent/run-xagent.sh
/opt/xagent/certs/ca-bundle-mega.crt
/etc/xagent/remGenie/remGenie.dev.env
/etc/xagent/bst/bst.dev.env
/etc/init.d/XAgent
