# netgear/RAX5/etc/hotplug.d/iface/99-checkip
[ "$ACTION" != "ifup" -a "$ACTION" != "ifupdate" -a "$ACTION" != "ifdown" ] && return

if [ "$INTERFACE" == "wan" ] || [ "$INTERFACE" = "pptp" ] || [ "$INTERFACE" = "l2tp" ]; then
    echo "check ip conflict" > /dev/console
    (/usr/sbin/check_ipconflict.sh "$INTERFACE" "$DEVICE" $ACTION)&
fi
