# tplink_en/Archer AX6000/etc/init.d/speedtest_status
#!/bin/sh /etc/rc.common

START=99
STOP=99

start() {
    speedtest_polling&
	return 0
}

stop() {
	return 0
}
