# tplink_en/Archer C60/etc/hotplug.d/button/45-wps-reset-led

do_blink()
{
    sleep 5 
	echo "do_blink" >/dev/console
	[ -f /tmp/reset_button ] && ledcli STATUS_BLINK
}

if [ "$ACTION" = "pressed" -a "$BUTTON" = "wps" ]; then
    [ ! -f /tmp/reset_button ] && echo "pressed" >/tmp/reset_button
#	do_blink &
fi
