# tplink_en/Archer C60/etc/hotplug.d/button/51-wps-reset
#
#
# Copyright (c) 2014, The Linux Foundation. All rights reserved.
#
#  Permission to use, copy, modify, and/or distribute this software for any
#  purpose with or without fee is hereby granted, provided that the above
#  copyright notice and this permission notice appear in all copies.
#
#  THE SOFTWARE IS PROVIDED "AS IS" AND THE AUTHOR DISCLAIMS ALL WARRANTIES
#  WITH REGARD TO THIS SOFTWARE INCLUDING ALL IMPLIED WARRANTIES OF
#  MERCHANTABILITY AND FITNESS. IN NO EVENT SHALL THE AUTHOR BE LIABLE FOR
#  ANY SPECIAL, DIRECT, INDIRECT, OR CONSEQUENTIAL DAMAGES OR ANY DAMAGES
#  WHATSOEVER RESULTING FROM LOSS OF USE, DATA OR PROFITS, WHETHER IN AN
#  ACTION OF CONTRACT, NEGLIGENCE OR OTHER TORTIOUS ACTION, ARISING OUT OF
#  OR IN CONNECTION WITH THE USE OR PERFORMANCE OF THIS SOFTWARE.
#


if [ "$ACTION" = "released" -a "$BUTTON" = "wps" ]; then
	if [ "$(uci get factory.factorymode.enable)" = "yes" -a "$(/sbin/is_cal)" != "true" ]; then
		touch /tmp/button_reset_check
	else
		if [ "$SEEN" -gt 3 ]; then
			echo "" > /dev/console
			echo "RESET TO FACTORY SETTING EVENT DETECTED" > /dev/console
			echo "PLEASE WAIT WHILE REBOOTING THE DEVICE..." > /dev/console
			#nvrammanager -e  -p  user-config  >/dev/null 2>&1
			#nvrammanager -r /tmp/reset-defconfig.xml -p  default-config  >/dev/null 2>&1
			#nvrammanager -w /tmp/reset-defconfig.xml -p  user-config   >/dev/null 2>&1	
			#rm -f /tmp/reset-defconfig.xml >/dev/null 2>&1
			#nvrammanager -e  -p  certificate  >/dev/null 2>&1	
			#rm -rf /overlay/*
			lua -e 'require("luci.controller.admin.firmware").hard_reset()'
		else
			ledcli STATUS_ON
		fi
	fi
fi

if [ "$ACTION" = "pressed" -a "$BUTTON" = "wps" ]; then
	if [ "$(uci get factory.factorymode.enable)" = "yes" -a "$(/sbin/is_cal)" != "true" ]; then
		echo "--> RESET button pressed" > /dev/console
	else
		sleep 4
		ledcli STATUS_BLINK
	fi
fi

