# netgear/RAX5/etc/hotplug.d/phy/00-ledstate
#echo $PORT $ACTION $SPEED > /tmp/hotplug_phy
# uci -c /var/state set hotplug_phy.port.num=$PORT
uci -c /var/state set hotplug_phy.port.action_$PORT=$ACTION
uci -c /var/state set hotplug_phy.port.speed_$PORT=$SPEED
uci -c /var/state commit hotplug_phy
killall -USR1 led_control
