# netgear/RAX5/etc/hotplug.d/phy/01-wanstate
if [ $PORT -eq 0 ]; then
	if [ "$ACTION" = linkup ]; then
		ifup wan
	else
		ifdown wan
	fi

fi

