# tplink_en/Archer AX90/usr/lib/opkg/info/family_time.list
/lib/family_time/core_family_time.sh
/sbin/family_check
/lib/family_time/config.sh
/sbin/bonusTime_check
/etc/init.d/family_time
/lib/family_time/core_log.sh
/lib/family_time/core.sh
