# tplink_en/Archer AX90/usr/lib/opkg/info/safethings.list
/usr/share/.avira/avirasentinelfull/avirasentinelfull_eng_build-2102018_prod/mainConfig.config.sig
/usr/bin/avira/aviraserviceselector/aviraserviceselector_eng_build-2102018_prod/aviraserviceselector
/usr/lib/avira/libcap.so.2.25
/usr/lib/avira/libnetfilter_conntrack.so.3
/usr/bin/avira/avirasentinellite/avirasentinellite_eng_build-2102018_prod/avirasentinellite.sh
/usr/share/.avira/avirasentinellite/avirasentinellite_eng_build-2102018_prod/passwords.txt
/usr/lib/avira/libnetfilter_conntrack.txt
/usr/bin/avira/aviraserviceselector/aviraserviceselector_eng_build-2102018_prod/uninstallSpecific.sh
/usr/bin/avira/aviraserviceselector/aviraserviceselector
/usr/share/.avira/avirasentinellite/avirasentinellite_eng_build-2102018_prod/mainConfig.config.sig
/usr/lib/avira/libwebsockets.txt
/usr/bin/avira/aviraserviceselector/aviraserviceselector_eng_build-2102018_prod/aviraserviceselector.sh
/usr/lib/avira/libnfnetlink.so.0
/usr/bin/avira/avirasentinelfull/avirasentinelfull_eng_build-2102018_prod/uninstallComplete.sh
/usr/bin/avira/avirasentinellite/avirasentinellite_eng_build-2102018_prod/uninstallSpecific.sh
/usr/share/.avira/avirasentinellite/avirasentinellite_eng_build-2102018_prod/ports.csv
/usr/share/.avira/avirasentinelfull/avirasentinelfull_eng_build-2102018_prod/ports.csv
/etc/init.d/aviraserviceselector
/usr/share/.avira/avirasentinelfull/avirasentinelfull_eng_build-2102018_prod/passwords.txt
/etc/init.d/avirasentinelfull
/usr/bin/avira/avirasentinellite/avirasentinellite_eng_build-2102018_prod/avirasentinellite
/etc/init.d/avirasentinellite
/usr/lib/avira/libcap.txt
/usr/share/.avira/avirasentinelfull/avirasentinelfull_eng_build-2102018_prod/mainConfig.config
/usr/lib/avira/libpcap.so.1
/usr/lib/avira/libcap.so.2
/usr/lib/avira/libmnl.so.0
/usr/bin/avira/avirasentinelfull/avirasentinelfull
/usr/bin/avira/sshpass
/usr/lib/avira/libnfnetlink.txt
/usr/sbin/aviraservice
/usr/lib/avira/libmnl.so.0.2.0
/usr/lib/avira/libwebsockets.so.13
/usr/share/.avira/aviraserviceselector/aviraserviceselector_eng_build-2102018_prod/mainConfig.config.sig
/usr/lib/avira/libnfnetlink.so.0.2.0
/usr/bin/avira/aviraserviceselector/aviraserviceselector_eng_build-2102018_prod/uninstallComplete.sh
/usr/bin/avira/avirasentinelfull/avirasentinelfull_eng_build-2102018_prod/uninstallSpecific.sh
/usr/lib/avira/libpcap.so.1.8.1
/usr/bin/avira/avirasentinellite/avirasentinellite
/usr/bin/avira/avirasentinelfull/avirasentinelfull_eng_build-2102018_prod/avirasentinelfull
/usr/bin/avira/avirasentinellite/avirasentinellite_eng_build-2102018_prod/uninstallComplete.sh
/usr/lib/avira/libmnl.txt
/usr/lib/avira/libnetfilter_conntrack.so.3.6.0
/usr/lib/avira/libpcap.txt
/usr/bin/avira/avirasentinelfull/avirasentinelfull_eng_build-2102018_prod/avirasentinelfull.sh
/usr/share/.avira/avirasentinellite/avirasentinellite_eng_build-2102018_prod/mainConfig.config
/usr/share/.avira/aviraserviceselector/aviraserviceselector_eng_build-2102018_prod/mainConfig.config
