# tplink_en/Archer AX90/usr/lib/opkg/info/url-class.list
/etc/init.d/url_class
/usr/sbin/url-class
