# tplink_en/Archer AX90/usr/lib/opkg/info/watchdog.list
/usr/lib/avira/libstdc++.so.6
/usr/lib/avira/libc-2.26.so
/usr/lib/avira/libstdc++.so.6.0.21
/usr/lib/avira/libz.txt
/usr/lib/avira/librt.txt
/etc/init.d/avirawatchdog
/usr/lib/avira/libmbedx509.so.2.16.7
/usr/lib/avira/libc.txt
/usr/lib/avira/libsodium.so.18
/usr/lib/avira/libz.so.1.2.11
/usr/lib/avira/libmbedtls.so.2.16.7
/usr/share/.avira/ca-certificates.config.sig
/usr/lib/avira/libsodium.txt
/usr/lib/avira/libm.so.6
/usr/lib/avira/libgcc_s.so.1
/usr/share/.avira/avirawatchdog/avirawatchdog_eng_build-2102018_prod/mainConfig.config
/usr/lib/avira/libz.so.1
/usr/lib/avira/libpthread-2.26.so
/usr/lib/avira/librt-2.26.so
/usr/bin/avira/avirawatchdog/avirawatchdog_eng_build-2102018_prod/avirawatchdog.sh
/usr/bin/avira/avirawatchdog/avirawatchdog_eng_build-2102018_prod/uninstallComplete.sh
/usr/bin/avira/avirawatchdog/avirawatchdog_eng_build-2102018_prod/avirawatchdog
/usr/lib/avira/libmbedtls.txt
/usr/lib/avira/libm.txt
/usr/lib/avira/libmbedcrypto.txt
/usr/lib/avira/libstdc++.txt
/usr/lib/avira/libpthread.so.0
/usr/lib/avira/libcurl.so.4.5.0
/usr/lib/avira/libm-2.26.so
/usr/bin/avira/avirawatchdog/avirawatchdog
/usr/scripts/avira/sendWatchdogHeartbeat.sh
/usr/lib/avira/libcurl.txt
/usr/lib/avira/libmbedcrypto.so.3
/usr/lib/avira/libmbedx509.so.0
/usr/lib/avira/libc.so.6
/usr/lib/avira/librt.so.1
/usr/lib/avira/libcurl.so.4
/usr/lib/avira/libsodium.so.18.3.0
/usr/lib/avira/libmbedx509.txt
/usr/share/.avira/avirawatchdog/avirawatchdog_eng_build-2102018_prod/mainConfig.config.sig
/usr/lib/avira/libpthread.txt
/usr/scripts/avira/disableWatchdogHeartbeat.sh
/usr/bin/avira/avirawatchdog/avirawatchdog_eng_build-2102018_prod/uninstallSpecific.sh
/usr/lib/avira/libmbedcrypto.so.2.16.7
/usr/share/.avira/ca-certificates.config
/usr/lib/avira/libgcc_s.txt
/usr/scripts/avira/watchdogCronTask.sh
/usr/scripts/avira/enableWatchdogHeartbeat.sh
/usr/lib/avira/libmbedtls.so.12
