# netgear/RAX36S/etc/init.d/DALH
#!/bin/sh /etc/rc.common
START=99

start() {
    /usr/bin/dalh 2>/dev/null &
}

stop() {
    killall dalh
    kill -17 1
    rm /tmp/dalh.log
}

boot() {
    start
}

restart() {
    stop
    start
}
