# netgear/RAX36S/etc/init.d/quagga
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=60

start() {
	/usr/sbin/quagga.init start $daemon_args
}

stop() {
	/usr/sbin/quagga.init stop $daemon_args
}

daemon_args=$@
