# tplink_en/Archer AX50/usr/lib/opkg/info/nvrammanager.list
/usr/sbin/rootfs-backup
/usr/bin/nvrammanager
/etc/init.d/rootfs_backup
