# tplink_en/Archer AX50/usr/lib/opkg/info/product-file.list
/etc/init.d/hnat-start
/usr/.gitignore
/lib/usb/reduction.sh
