# tplink_en/Archer AXE75/usr/lib/opkg/info/safethings.list
/usr/bin/avira/avirasentinellite/avirasentinellite_rel-210604-2_prod/uninstallSpecific.sh
/usr/share/.avira/avirasentinellite/avirasentinellite_rel-210604-2_prod/ports.csv
/usr/share/.avira/aviraserviceselector/aviraserviceselector_rel-210604-2_prod/mainConfig.config
/usr/bin/avira/aviraserviceselector/aviraserviceselector_rel-210604-2_prod/aviraserviceselector
/usr/bin/avira/aviraserviceselector/aviraserviceselector_rel-210604-2_prod/aviraserviceselector.sh
/usr/bin/avira/avirasentinelfull/avirasentinelfull_rel-210604-2_prod/avirasentinelfull.sh
/usr/lib/avira/libcap.so.2.27
/usr/lib/avira/libnetfilter_conntrack.so.3
/usr/share/.avira/avirasentinelfull/avirasentinelfull_rel-210604-2_prod/mainConfig.config.sig
/usr/lib/avira/libnetfilter_conntrack.txt
/usr/bin/avira/avirasentinelfull/avirasentinelfull_rel-210604-2_prod/avirasentinelfull
/usr/bin/avira/aviraserviceselector/aviraserviceselector
/usr/bin/avira/aviraserviceselector/aviraserviceselector_rel-210604-2_prod/uninstallSpecific.sh
/usr/lib/avira/libwebsockets.txt
/usr/lib/avira/libnfnetlink.so.0
/usr/bin/avira/avirasentinellite/avirasentinellite_rel-210604-2_prod/avirasentinellite
/etc/init.d/aviraserviceselector
/etc/init.d/avirasentinelfull
/etc/init.d/aviraservicemaster
/usr/bin/avira/avirasentinelfull/avirasentinelfull_rel-210604-2_prod/uninstallComplete.sh
/etc/init.d/avirasentinellite
/usr/lib/avira/libcap.txt
/usr/share/.avira/avirasentinelfull/avirasentinelfull_rel-210604-2_prod/ports.csv
/usr/share/.avira/avirasentinelfull/avirasentinelfull_rel-210604-2_prod/passwords.txt
/usr/lib/avira/libpcap.so.1
/usr/lib/avira/libcap.so.2
/usr/lib/avira/libmnl.so.0
/usr/bin/avira/avirasentinellite/avirasentinellite_rel-210604-2_prod/uninstallComplete.sh
/usr/bin/avira/avirasentinelfull/avirasentinelfull
/usr/share/.avira/avirasentinelfull/avirasentinelfull_rel-210604-2_prod/mainConfig.config
/usr/bin/avira/sshpass
/usr/lib/avira/libnfnetlink.txt
/usr/sbin/aviraservice
/usr/share/.avira/avirasentinellite/avirasentinellite_rel-210604-2_prod/mainConfig.config
/usr/lib/avira/libmnl.so.0.2.0
/usr/lib/avira/libwebsockets.so.13
/usr/lib/avira/libnfnetlink.so.0.2.0
/usr/bin/avira/avirasentinellite/avirasentinellite_rel-210604-2_prod/avirasentinellite.sh
/usr/share/.avira/avirasentinellite/avirasentinellite_rel-210604-2_prod/passwords.txt
/usr/share/.avira/aviraserviceselector/aviraserviceselector_rel-210604-2_prod/mainConfig.config.sig
/usr/bin/avira/aviraserviceselector/aviraserviceselector_rel-210604-2_prod/uninstallComplete.sh
/usr/bin/avira/avirasentinellite/avirasentinellite
/usr/bin/avira/avirasentinelfull/avirasentinelfull_rel-210604-2_prod/uninstallSpecific.sh
/usr/lib/avira/libmnl.txt
/usr/lib/avira/libpcap.so.1.9.1
/usr/share/.avira/avirasentinellite/avirasentinellite_rel-210604-2_prod/mainConfig.config.sig
/usr/lib/avira/libpcap.txt
/usr/lib/avira/libnetfilter_conntrack.so.3.7.0
