# tplink_en/Archer C2600/usr/lib/opkg/info/ledctrl.list
/etc/init.d/ledctrl
/usr/bin/ledcli
/usr/bin/ledctrl
