# tplink_en/Archer C2600/usr/lib/opkg/info/openslp.list
/etc/slp.spi
/usr/sbin/slpd
/usr/lib/libslp.so.1.0.1
/etc/init.d/slpd
/usr/lib/libslp.so.1
/etc/slp.reg
/lib/upgrade/keep.d/openslp
/usr/bin/slptool
/etc/slp.conf
