# tplink_en/Archer C2600/usr/lib/opkg/info/switch-setup.list
/lib/switch/config.sh
/sbin/switch_link_down
/sbin/switch_bind_mac
/lib/switch/core.sh
/etc/init.d/switch
/etc/init.d/phyport
/etc/init.d/downports
/sbin/switch_link_up
/etc/config/switch
