# netgear/RAXE500/usr/bin/circlev2/shares/usr/bin/service/timetracker/run
#!/bin/sh
echo "******timetracker run script******"
killall -9 timetracker
sleep 1
exec $CIRCLE_ROOT/timetracker
