# netgear/RS700S/usr/bin/circlev2/shares/usr/bin/service/webd/run
#!/bin/sh
exec $CIRCLE_ROOT/webd -u root -E $CIRCLE_DATA/webd.pem  -D
