# tplink_en/AD7200/usr/lib/opkg/info/client_mgmt.list
/etc/hotplug.d/iface/26-client_mgmt
/usr/bin/client_mgmt
/etc/init.d/client_mgmt
