# tplink_en/AD7200/usr/lib/opkg/info/kmod-qca-wifi-10.4-akronite-perf.list
/lib/wifi/tplink.sh
/lib/modules/3.4.103/umac.ko
/lib/modules/3.4.103/ath_hal.ko
/etc/hotplug.d/button/40-wireless
/usr/sbin/thermaltool
/usr/sbin/athadhoc
/lib/wifi/tplink-update-uci
/lib/modules/3.4.103/ath_rate_atheros.ko
/usr/sbin/apstats
/usr/sbin/spectraltool
/lib/modules/3.4.103/smart_antenna.ko
/usr/sbin/pktlogdump
/usr/sbin/80211stats
/lib/modules/3.4.103/asf.ko
/etc/init.d/wireless
/usr/sbin/radartool
/lib/wifi/qcawifi_countrycode.txt
/lib/modules/3.4.103/adf.ko
/usr/sbin/wifitool
/usr/lib/libtlvparser.so
/usr/lib/libtlvencoder.so
/usr/sbin/athtestcmd
/usr/sbin/tx99tool
/lib/modules/3.4.103/hst_tx99.ko
/usr/sbin/athstats
/usr/lib/libtlvtemplate.so
/etc/hotplug.d/firewall/02-wireless-vlan
/usr/lib/libtlvcmdrsp.so
/lib/modules/3.4.103/ath_dfs.ko
/sbin/wifi_try
/usr/lib/libtlvutil.so
/usr/sbin/wlanconfig
/lib/modules/3.4.103/ath_pktlog.ko
/usr/sbin/athdiag
/lib/modules/3.4.103/ath_spectral.ko
/lib/modules/3.4.103/ath_dev.ko
/usr/sbin/athstatsclr
/usr/sbin/pktlogconf
/sbin/is_cal
/etc/hotplug.d/iface/02-wireless-vlan
/lib/wifi/qcawifi.sh
/etc/modules.d/33-qca-wifi-10.4-akronite-perf
