# tplink_en/AD7200/usr/lib/opkg/info/ledpm.list
/lib/ledpm/ledpm_core.sh
/sbin/ledpm_onoff
/sbin/ledpm
/sbin/ledgeneral
/lib/ledpm/ledgeneral_core.sh
/etc/config/ledpm
/sbin/ledgeneral_onoff
/etc/init.d/ledgeneral
/lib/ledpm/ledpm_function.sh
/etc/init.d/ledpm
