# tplink_en/Archer C58HP/etc/init.d/client_mgmt
#!/bin/sh /etc/rc.common

SERVICE_DAEMONIZE=1
#SERVICE_WRITE_PID=1

START=26
STOP=26

start()
{
	service_start /usr/bin/client_mgmt
	return 0
}

stop()
{
	service_stop /usr/bin/client_mgmt
	return 0
}
