# tplink_en/Archer C5400X/usr/lib/opkg/info/kmod-geforce-now.list
/etc/init.d/gfn
/lib/modules/iplatform/gfn.ko
/etc/modules.d/91-geforce-now
