# tplink_en/Archer C5400X/usr/lib/opkg/info/monit.list
/usr/sbin/monit
/etc/init.d/zzzmonit
/lib/monit/monit_proclist
