# tplink_en/Archer C5400X/usr/lib/opkg/info/mynginx.list
/etc/nginx/nginx.conf
/etc/nginx/fastcgi_params
/etc/nginx/koi-win
/usr/lib/lua/luci/model/tpFileVer.lua
/usr/sbin/nginx
/etc/nginx/win-utf
/etc/certificate/cert.pem
/etc/nginx/koi-utf
/etc/init.d/nginx
/etc/nginx/mime.types
/lib/upgrade/keep.d/mynginx
/www/html/50x.html
/etc/certificate/key.pem
/www/html/index.html
/usr/sbin/create_cer.sh
/usr/sbin/internet_operation.sh
