# tplink_en/Archer C5400X/usr/lib/opkg/info/p2pd.list
/usr/bin/upnpc
/usr/bin/stuncli
/usr/bin/p2pd
/etc/init.d/p2pd
/etc/hotplug.d/iface/.gitignore
/usr/.gitignore
/usr/bin/p2ptest
