# tplink_en/Archer C5400X/usr/lib/opkg/info/print-server.list
/lib/print_server/printer_discovery.sh
/etc/init.d/print_server
/etc/hotplug.d/usb/11-printer
/lib/upgrade/keep.d/print-server
/etc/hotplug.d/iface/91-printer
