# tplink_en/Archer C5400X/etc/hotplug.d/iface/91-printer
[ "$INTERFACE" == "lan" -a "$ACTION" == "ifup" ] && {
    /etc/init.d/print_server start
}

[ "$INTERFACE" == "lan" -a "$ACTION" == "ifdown" ] && {
    /etc/init.d/print_server stop
}
