# ubiquiti/UAP-AC-PRO/etc/init.d/ubnt_sysinit
#!/bin/sh /etc/rc.common

START=00

boot() {
	# disable write back caching
	[ -f /usr/etc/rc.d/rc.disable_wb_cache ] && . /usr/etc/rc.d/rc.disable_wb_cache

	. /usr/etc/rc.d/rc.sysinit
}
