# intelbras/WRH211/etc/inittab
:0:sysinit:/bin/mount -a
#si::sysinit:/etc/rc.d/rc.sysinit
:0:sysinit:/etc/rc.d/rc.bridge
7:2345:respawn:/sbin/getty -L ttyS 9600
#::respawn:/sbin/getty -L ttyS 9600 
::askfirst:/bin/ash
#::askfirst:/bin/ash

